# EmbeddedLinux_Kitchen/Task7/SW_ROOT/boot/rootramfs/etc/inittab
::sysinit:/etc/init.d/rcS
console::sysinit:/bin/sh
